# build_files/system_files/usr/lib/systemd/system/lvs-firstboot-setup.service
[Unit]
Description=Linux Validation Suite first boot setup
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lvs-firstboot-setup

[Install]
WantedBy=multi-user.target
