# build_files/system_files/etc/skel/.config/autostart/lvs-login-preflight.desktop
[Desktop Entry]
Type=Application
Name=LVS Login Preflight
Comment=Checks network, time sync, and update status before validation
Exec=konsole --hold -e /usr/libexec/lvs/lvs-login-preflight
Terminal=false
X-KDE-autostart-after=panel
