# build_files/system_files/etc/skel/.config/autostart/org.kde.plasma-welcome.desktop
[Desktop Entry]
Type=Application
Name=Plasma Welcome
Exec=plasma-welcome
Hidden=true
